# massbuild-status.timer
[Unit]
Description=Update the OpenMandriva mass build status report

[Timer]
OnUnitActiveSec=1h
AccuracySec=5min
Persistent=true
OnBootSec=1min

[Install]
WantedBy=timers.target
